# .docker/php/jobs.ini
[program:read-url-from-queue]
command=php bin/read-url-from-queue
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
priority=100
user=www-data

[program:send-urls-to-queue]
command=php bin/send-urls-to-queue
numprocs=1
autostart=true
autorestart=false
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
priority=200
user=www-data